# com.dev-log-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dev-log-monitor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/dd/Desktop/1_dev/coding-rule2/projects/post_tool/venv/bin/python</string>
        <string>/Users/dd/Desktop/1_dev/coding-rule2/projects/post_tool/automation/dev_log_watcher.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/dd/Desktop/1_dev/coding-rule2/projects/post_tool</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/dd/Desktop/1_dev/coding-rule2/projects/post_tool/logs/dev_monitor.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/dd/Desktop/1_dev/coding-rule2/projects/post_tool/logs/dev_monitor_error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>ANTHROPIC_API_KEY</key>
        <string>your_api_key_here</string>
    </dict>
</dict>
</plist>